# local_packages/moos-special-keys/moos-special-keys.install
post_install() {
    systemctl enable special-keys-backlight
}

pre_remove() {
    systemctl disable special-keys-backlight
}
